# .config/ugly-hacks/rules/80-hotplug-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", HOTPLUG=="1", RUN+="/bin/systemctl start hotplug-monitor"
